# powerhausbox-core/rootfs/etc/profile.d/powerhausbox-host-tools.sh
export PATH="/usr/local/bin:${PATH}"

cat <<'EOF'

PowerHausBox host-capable tools:
  host-shell   Enter the Home Assistant host namespaces
  ha-host ...  Run the host HA CLI from this add-on shell

EOF
